# scripts/com.kaipanla.ocr.service.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kaipanla.ocr.service</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/macclaw/.openclaw/workspace/scripts/sentiment_ocr_service.py</string>
        <string>--server</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/macclaw/.openclaw/workspace</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/tmp/kaipanla/ocr_service_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/kaipanla/ocr_service_stderr.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
